# install/com.bsd.salescopilot.env.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bsd.salescopilot.env</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/launchctl</string>
        <string>setenv</string>
        <string>BSD_COPILOT_PATH</string>
        <string>__BSD_COPILOT_PATH__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
